# overlay/etc/local.d/10-lte_setup.start
INET="eth1"
while [ ! -e /sys/class/net/$INET ]
do
  sleep 2
done
echo Configuring external network conection.
ifup $INET
