# overlay/etc/local.d/20-set-passwords.start
root_pw=$(pwgen -cnsB 12 1)
wlan_pw=$(pwgen -cnsB 12 1)

echo Root PW: $root_pw
echo Wlan PW: $wlan_pw
